# conf/supervisor_celery_flower.ini
[program:imooc-celery-flower]
command=python manage.py celery flower
directory=/Users/pinghua.wph/Programme/PycharmProjects/celery-learning/my_django_celery
environment=PATH="/opt/anaconda2/envs/python36_celery/bin"
stdout_logfile=/Users/pinghua.wph/Programme/PycharmProjects/celery-learning/my_django_celery/logs/celery.flower.log
stderr_logfile=/Users/pinghua.wph/Programme/PycharmProjects/celery-learning/my_django_celery/logs/celery.flower.log
autostart=true
autorestart=true
startsecs=10
stopwatisecs=60
priority=995
