# conf/supervisor_celery_worker_2.ini
[program:imooc-celery-worker-2]
command=python manage.py celery worker -l INFO
directory=/Users/pinghua.wph/Programme/PycharmProjects/celery-learning/my_django_celery
environment=PATH="/opt/anaconda2/envs/python36_celery/bin"
stdout_logfile=/Users/pinghua.wph/Programme/PycharmProjects/celery-learning/my_django_celery/logs/celery.worker2.log
stderr_logfile=/Users/pinghua.wph/Programme/PycharmProjects/celery-learning/my_django_celery/logs/celery.worker2.log
autostart=true
autorestart=true
startsecs=10
stopwatisecs=60
priority=993
