# conf/supervisor_celery_django.ini
[program:imooc-celery-django]
command=python manage.py runserver 0.0.0.0:8000
directory=/Users/pinghua.wph/Programme/PycharmProjects/my_django_celery
environment=PATH="/opt/anaconda2/envs/python36_celery/bin"
stdout_logfile=/Users/pinghua.wph/Programme/PycharmProjects/my_django_celery/logs/celery.django.log
stderr_logfile=/Users/pinghua.wph/Programme/PycharmProjects/my_django_celery/logs/celery.django.log
autostart=true
autorestart=true
startsecs=10
stopwatisecs=60
priority=994
