# conf/supervisor_celery_redis.ini
[program:celery-redis]
command=/usr/local/bin/redis-server /usr/local/etc/redis.conf
environment=PATH="/opt/anaconda2/envs/python36_celery/bin"
directory=/Users/pinghua.wph/Programme/PycharmProjects/my_django_celery
stdout_logfile=/Users/pinghua.wph/Programme/PycharmProjects/my_django_celery/logs/celery.redis.log
stderr_logfile=/Users/pinghua.wph/Programme/PycharmProjects/my_django_celery/logs/celery.redis.log
autostart=true
autorestart=true
