# conf/supervisor_celery_beat.ini
[program:imooc-celery-beat]
command=python manage.py celery beat -l INFO
directory=/Users/pinghua.wph/Programme/PycharmProjects/celery-learning/my_django_celery
environment=PATH="/opt/anaconda2/envs/python36_celery/bin"
stdout_logfile=/Users/pinghua.wph/Programme/PycharmProjects/celery-learning/my_django_celery/logs/celery.beat.log
stderr_logfile=/Users/pinghua.wph/Programme/PycharmProjects/celery-learning/my_django_celery/logs/celery.beat.log
autostart=true
autorestart=true
startsecs=10
stopwatisecs=60
priority=998
